# declaration-management.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/warehouse_management.sock

[Install]
WantedBy=sockets.target
